# private_dot_config/autostart/numlock.desktop
[Desktop Entry]
Name=Numlockx
Exec=/usr/bin/numlockx on
Type=Application
Categories=Utility
