# private_dot_config/autostart/remaps.desktop
[Desktop Entry]
Name=Remaps
Exec=/home/david/.local/bin/remaps
Type=Application
Categories=Utility
